# _supervisors/redis.conf
[program:redis-server-worker]
process_name=%(program_name)s_%(process_num)02d
command=redis-server
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/path/to/project/webroot/storage/logs/_redis-server-worker.log